# setup/server/zai-proxy.service
[Unit]
Description=Z.ai 1M-first Proxy (conditional short fallback)
After=network.target
# cokacdir가 Claude Code를 spawn하기 전에 프록시가 먼저 올라와야 함
Before=cokacdir.service

[Service]
Type=simple
ExecStart=/usr/bin/node /home/ubuntu/scripts/zai-fallback-proxy.js
Restart=always
RestartSec=3
Environment=NODE_ENV=production
# 튜닝 (필요시 변경)
# Environment=ZAI_UPSTREAM=https://api.z.ai/api/anthropic
# Environment=PROXY_PORT=8788
# 조건부 fallback: 짧은 요청만 표준 모델로 살리고, 긴 요청은 1M을 보장한다.
# 짧은 복구/부트 요청은 glm-5.2 계열 과부하 시 5.1 → 4.6 순서로 회피한다.
# 긴 요청은 프록시 코드가 reqBytes 기준으로 fallback을 막아 1M/긴문맥 보장을 우선한다.
Environment=FALLBACK_MODELS=glm-5.1,glm-4.6
Environment=OVERLOAD_MODELS=glm-5.2[1m],glm-5.2
Environment=STANDARD_FALLBACK_MAX_BYTES=307200
Environment=NOTIFY_CHAT=-5402233196
Environment=NOTIFY_ALL_MODELS=1
Environment=CONTEXT_LIMIT_TOKENS=1000000
Environment=CONTEXT_WARN_PCT=70
Environment=CONTEXT_DANGER_PCT=85
Environment=MAX_ATTEMPTS=4
# Environment=NOTIFY_BOT_TOKEN=  (비워두면 cokacctl.json 첫 토큰 자동 사용)
# Environment=MAX_ATTEMPTS=4
StandardOutput=append:/home/ubuntu/.local/state/cokacdir/zai-fallback-proxy.stdout.log
StandardError=append:/home/ubuntu/.local/state/cokacdir/zai-fallback-proxy.log

[Install]
WantedBy=default.target
